# recipes-bsp/gpuconf/gpuconf/78-tegra-gpu.rules
ACTION=="add", KERNEL=="renderD128", RUN+="/usr/bin/tegra-gpu"
